# supervisor/nbserver.conf
[program:nbserver]
user = mvd
directory = /home/mvd
command = /home/mvd/miniconda3/bin/ipython notebook --profile=nbserver
redirect_stderr = true
stdout_logfile = /var/log/nbserver.log
